# delegator/aul-delegator-server.service
[Unit]
Description=Delegator

[Service]
BusName=org.tizen.aul.delegator
ExecStart=/usr/bin/aul-delegator-server

